# discord-voice-bot.service
[Unit]
Description=Discord Voice Bot
After=network-online.target
Wants=network-online.target
# Restart bot if it crashes, with delay to avoid rapid restart loops

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/discord-voice-bot
ExecStart=/home/ubuntu/discord-voice-bot/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits — conservative for Free Tier ARM
MemoryMax=256M
CPUQuota=25%

# Environment hardening (optional but good practice)
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target